# files/root/etc/init.d/systemup
#!/bin/sh /etc/rc.common
START=80

start()
{
	/usr/sbin/systemup.sh --svc &
}

stop()
{
	ps | grep systemup.sh | grep -v 'grep' | awk '{print $1}' | xargs kill
}

run_reboot()
{
	stop
	start
}
